# service/ax25node.service
[Unit]
Description=AX.25 Node

[Service]
ExecStart=/home/kd9yqk/ham314/launchers/ax25service.sh
WorkingDirectory=/home/kd9yqk/ham314/launchers/
StandardOutput=inherit
StandardError=inherit
Restart=no
User=kd9yqk
TimeoutStopSec=3
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
